# services/joyosc.service
[Unit]
Description=joyosc

[Service]
ExecStart=/usr/local/bin/joyosc

[Install]
WantedBy=default.target
